# debian/systemd/steam-session.service
[Unit]
Description=Steam
PartOf=graphical-session.target

[Service]
ExecStart=/usr/bin/steam-session

